# go1.25/rootfs/etc/s6-overlay/s6-rc.d/goapp/run
#!/command/with-contenv bash

set -eou pipefail

cd /app

exec s6-setuidgid goapp /app/binary
